# os/build_files/config/systemd/apparatus-first-login.service
[Unit]
Description=Apparatus First Login Setup
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/apparatus/first-login.sh
RemainAfterExit=no

[Install]
WantedBy=graphical-session.target
