# os/build_files/config/systemd/apparatus-firstboot.service
[Unit]
Description=Apparatus First Boot Setup
ConditionFirstBoot=yes
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/apparatus/firstboot-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
