# os/build_files/config/systemd/apparatus-firstboot.service
[Unit]
Description=Apparatus First Boot Setup
ConditionPathExists=!/var/lib/apparatus/firstboot-done
After=local-fs.target
Before=gdm.service display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/apparatus/firstboot-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
